# ops/systemd/kucoin-monitoring-daily-analysis.service
[Unit]
Description=KuCoin lane daily monitoring analysis
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=we4free
WorkingDirectory=/home/we4free/agent/repos/kucoin-lane
ExecStart=/usr/bin/env python3 /home/we4free/agent/repos/kucoin-lane/scripts/monitoring_automation.py analyze --period daily
