# ops/systemd/user/kucoin-monitoring-hourly.timer
[Unit]
Description=Run KuCoin lane hourly monitoring snapshot (user)

[Timer]
OnCalendar=hourly
Persistent=true
Unit=kucoin-monitoring-hourly.service

[Install]
WantedBy=timers.target
